# release/alert-bundle/grafana-anomaly-exporter.service
[Unit]
Description=Grafana Anomaly Exporter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=grafana-anomaly
Group=grafana-anomaly
WorkingDirectory=/opt/grafana-anomaly-exporter
EnvironmentFile=/etc/grafana-anomaly-exporter/exporter.env
ExecStart=/opt/grafana-anomaly-exporter/.venv/bin/python /opt/grafana-anomaly-exporter/main.py
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/grafana-anomaly-exporter /var/log/grafana-anomaly-exporter
StandardOutput=append:/var/log/grafana-anomaly-exporter/exporter.log
StandardError=append:/var/log/grafana-anomaly-exporter/exporter.log

[Install]
WantedBy=multi-user.target
